# setup/testdevice.service
[Unit]
Description=Wifi-Test-Device

[Service]
ExecStart=/home/pi/go/bin/wifi-test-device

[Install]
WantedBy=multi-user.target